# config/kalipi-agent.service
[Unit]
Description=KaliPi Agent — Remote monitoring API over Tailscale
After=network-online.target tailscaled.service
Wants=network-online.target
Requires=tailscaled.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/python3 -m agent
WorkingDirectory=/opt/kalipi
Environment=PYTHONPATH=/opt/kalipi
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
