# config/kalipi-monitor.service
[Unit]
Description=KaliPi LCD Security Monitor
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/kalipi/scripts/monitor-lcd.sh --loop
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
